# crontab
# Run automation every min
* * * * * bash /app/runSchedule.sh >> /var/log/myjob.log 2>&1
