# crontab
# Run automation every min
*/15 * * * * bash /app/runSchedule.sh >> /var/log/myjob.log 2>&1
